# etc/gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/github/cashhaikya
command=/home/ubuntu/github/cashhaikya/venv/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/github/cashhaikya/chk.sock cashhaikya.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log
[group:guni]
programs:gunicorn
